# mnt/var/lib/dpkg/info/cgroupfs-mount.list
/.
/etc
/etc/init.d
/etc/init.d/cgroupfs-mount
/lib
/lib/systemd
/lib/systemd/system
/usr
/usr/bin
/usr/bin/cgroupfs-mount
/usr/bin/cgroupfs-umount
/usr/share
/usr/share/doc
/usr/share/doc/cgroupfs-mount
/usr/share/doc/cgroupfs-mount/README.md
/usr/share/doc/cgroupfs-mount/changelog.gz
/usr/share/doc/cgroupfs-mount/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/cgroupfs-mount.8.gz
/lib/systemd/system/cgroupfs-mount.service
/usr/share/man/man8/cgroupfs-umount.8.gz
