# mnt/var/lib/dpkg/info/cgroupfs-mount.prerm
#!/bin/sh
set -e
# Automatically added by dh_installdeb
dpkg-maintscript-helper rm_conffile /etc/init/cgroupfs-mount.conf 1.4~ -- "$@"
# End automatically added section
# Automatically added by dh_installinit
if [ -x "/etc/init.d/cgroupfs-mount" ]; then
	invoke-rc.d cgroupfs-mount stop || exit $?
fi
# End automatically added section
