# mnt/var/lib/dpkg/info/redis-server.conffiles
/etc/default/redis-server
/etc/init.d/redis-server
/etc/logrotate.d/redis-server
/etc/redis/redis.conf
