# mnt/var/lib/dpkg/info/redis-server.list
/.
/etc
/etc/default
/etc/default/redis-server
/etc/init.d
/etc/init.d/redis-server
/etc/logrotate.d
/etc/logrotate.d/redis-server
/etc/redis
/etc/redis/redis.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/redis-server.service
/lib/systemd/system/redis-server@.service
/usr
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/redis-server
/usr/share/doc/redis-server/MANIFESTO.gz
/usr/share/doc/redis-server/README.md.gz
/usr/share/doc/redis-server/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/redis-server.1.gz
/usr/bin/redis-server
/usr/share/doc/redis-server/00-RELEASENOTES.gz
/usr/share/doc/redis-server/NEWS.Debian.gz
/usr/share/doc/redis-server/changelog.Debian.gz
